# smartpen/system/buildroot-intermediate/target/etc/profile.d/debugger.sh
if [ -f /usr/lib/libdebugger.so ];then
		 export LD_PRELOAD=/usr/lib/libdebugger.so
fi
